# init.d/yamaha2rest
### BEGIN INIT INFO
# Provides:          yamaha2rest
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Yamaha-2-REST gateway
# Description:       Service which provides a REST API to the RS-232 interface used by Yamaha RX-V receiver series
### END INIT INFO

PATH=/bin:/usr/bin:/sbin:/usr/sbin
DESC="Yamaha-2-REST gateway"
NAME=yamaha2rest
DAEMON=/opt/"$NAME"/bin/server.py
PIDFILE=/var/run/"$NAME".pid
SCRIPTNAME=/etc/init.d/"$NAME"

test -f $DAEMON || exit 0

. /lib/lsb/init-functions

[ -r /etc/default/"$NAME" ] && . /etc/default/"$NAME"

case "$1" in
start)
	log_daemon_msg "Starting $DESC" "$NAME"
	start_daemon -p $PIDFILE $DAEMON $EXTRA_OPTS
	log_end_msg $?
	;;
stop)
	log_daemon_msg "Stopping $DESC" "$NAME"
	killproc -p $PIDFILE $DAEMON
	RETVAL=$?
	[ $RETVAL -eq 0 ] && [ -e "$PIDFILE" ] && rm -f $PIDFILE
	log_end_msg $RETVAL
	;;
restart)
	log_daemon_msg "Restarting $DESC" "$NAME"
	$0 stop
	$0 start
	;;
reload|force-reload)
	log_daemon_msg "Reloading configuration files for $DESC" "$NAME"
	# Since we have no config, do nothing
	log_end_msg 0
	;;
status)
	status_of_proc -p $PIDFILE $DAEMON $NAME && exit 0 || exit $?
	;;
*)
	log_action_msg "Usage: $0 {start|stop|status|restart|reload|force-reload}"
	exit 2
	;;
esac
exit 0

